# com.gpxl.naenae.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.gpxl.naenae</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>-m</string>
    <string>naenae.app</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/gerlando/Documents/GitHub/watcher</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/gerlando/.naenae/logs/launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/gerlando/.naenae/logs/launchd.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/gerlando/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/gerlando</string>
    <key>NAENAE_HOME</key>
    <string>/Users/gerlando/.naenae</string>
  </dict>
</dict>
</plist>
